# RACCOLTA/STUNNEL/5.17/opt/Stunnel/5.17/Debian7/Linux/amd64/share/doc/stunnel/examples/stunnel.service
[Unit]
Description=SSL tunnel for network daemons
After=syslog.target

[Service]
ExecStart=/opt/Stunnel/5.17/Debian7/Linux/amd64/bin/stunnel
Type=forking

[Install]
WantedBy=multi-user.target
